# debian-bo-master/bo/source/sound/nas-1.2p4/debian.postrm
#! /bin/sh
if [ $1 = "purge" ]
then
  update-rc.d nas remove >/dev/null
fi
