# debian-bo-master/bo/source/sound/nas-1.2p4/debian.postinst
#! /bin/sh
update-rc.d nas defaults >/dev/null
