# debian-bo-master/bo/source/net/net-acct-0.4/debian.postrm
#! /bin/sh

set -e

if [ "$1" = "purge" ]
then
  update-rc.d net-acct remove >/dev/null
  echo "Please remove /var/log/net-acct/ manually!"
fi
