# debian-bo-master/bo/source/net/net-acct-0.4/debian.conffiles
/etc/naccttab
/etc/init.d/net-acct
/etc/cron.daily/net-acct
