# debian-bo-master/bo/source/utils/debmake-3.3.1/lib/init.postinst
# Automatically inserted by debstd on #DATE#
update-rc.d #SCRIPT# #INITPARAMS# >/dev/null
/etc/init.d/#SCRIPT# start
