# debian-bo-master/bo/source/mail/smail_3.2/samples/bigsite/bargw/news/crontab
#ident "@(#) crontab,v 1.2 1990/10/24 05:20:39 tron Exp"
# Entries for Smail and smail map generation, in root's crontab file
# on bargw.

#
# the following lines check the Smail spool directories for errors, and
# perform miscellaneous cleanup operations.
#
20  1 * * *   /usr/lib/smail/checkerr
40  1 * * *   /usr/lib/smail/savelog /usr/spool/smail/log/logfile
40  1 * * 1   /usr/lib/smail/savelog /usr/spool/smail/log/paniclog

#
# the following command unbatches USENET map entries once per day
 0  2 * * *   /usr/lib/smail/daily
