# debian-bo-master/bo/source/admin/file-rc-0.4.1/debian/postinst
#!/bin/sh
set -e
if [ -d /etc/rc0.d -a ! -f /etc/runlevel.conf ]; then
	/usr/lib/file-rc/rclink2file.sh > /etc/runlevel.conf
	cp /etc/runlevel.conf /etc/runlevel.fallback
	( cd / && tar czf /var/backups/rc-links.tar.gz etc/rc?.d )
	rm -rf /etc/rc?.d
fi
dpkg-divert --package file-rc --add --rename \
	--divert /usr/sbin/update-rc.d.links /usr/sbin/update-rc.d
cp /usr/lib/file-rc/update-rc.d /usr/sbin/update-rc.d
dpkg-divert --package file-rc --add --rename \
	--divert /etc/init.d/rc.links /etc/init.d/rc
cp /usr/lib/file-rc/rc /etc/init.d/rc
