# debian-bo-master/bo/source/admin/file-rc-0.4.1/debian/postrm
#!/bin/sh
set -e
if [ remove = "$1" ]; then
	mv /usr/sbin/update-rc.d.links /usr/sbin/update-rc.d
	dpkg-divert --package file-rc --remove \
		--divert /usr/sbin/update-rc.d.links /usr/sbin/update-rc.d
	mv -f /etc/init.d/rc.links /etc/init.d/rc
	dpkg-divert --package file-rc --remove \
		--divert /etc/init.d/rc.links /etc/init.d/rc
fi
if [ purge = "$1" ]; then
	rm -f /var/backups/rc-links.tar.gz /etc/runlevel.conf /etc/runlevel.fallback
	echo Please purge /var/backups/rc.* manually.
fi
