# debian-bo-master/bo/source/base/ld.so-1.8.10/util/ld.so.preload
	libm.so.5	
libX11.so.6
# bye-bye # howdy
libtcl7.4.so.1
