# debian-bo-master/bo/source/mail/deliver-2.1.12/debian.cron
#!/bin/sh
#
# cron script to rotate deliver log files
#
# Robert Leslie <rob@mars.org>

cd /var/log
for log in deliver.log deliver.errlog
do
  if [ -f $log ]
  then
    savelog -g adm -m 644 -u root -c 4 $log >/dev/null
  fi
done
