# debian-bo-master/bo/source/net/ipx-1.0/debian.postinst
#!/bin/sh
set -e

update-rc.d ipx defaults >/dev/null

/etc/init.d/ipx start
