# debian-bo-master/bo/source/utils/debmake-3.3.1/sample.multi/cron.daily
#!/bin/sh

cd /var/log/mgetty
for LOG in mg_*log; do
   if [ -f $LOG ]; then
      savelog -g adm -m 644 -u root -c 7 $LOG >/dev/null
   fi
done
