# debian-bo-master/bo/source/mail/sendmail-8.8.5/debian/postrm
#!/bin/sh -e
#
# Debian package postrm
# Version 1.2
#
# Robert Leslie <rob@mars.org>

case "$1" in
    remove)
    ;;

    upgrade)
    ;;

    purge)
	rm -rf /var/spool/mqueue
	rm -f /etc/sendmail.cf /etc/aliases.{db,dir,pag}
	rm -f /etc/mail/sendmail.{cf,cw,ct,mc}{,.old}
	rmdir /etc/mail 2>/dev/null || true
	update-rc.d sendmail remove >/dev/null
    ;;

    failed-upgrade|abort-install|abort-upgrade|disappear)
    ;;

    *)
	echo "postrm called with unknown argument \`$1'" >&2
	exit 0
    ;;
esac
