# debian-bo-master/bo/source/mail/smail_3.2/debian.conffiles
/etc/smail/crontab
/etc/cron.daily/smail
