# debian-bo-master/bo/source/mail/smail_3.2/debian.crontab
# Default `mail' user's crontab, supplied Smail for Debian Linux.
# See /usr/doc/smail/copyright for details, or read smail(5) for details of
# how to configure Smail.  Autoconfiguration done by /usr/sbin/smailconfig.
# (min, hour)	(dm,m,dw)	(command)
3,23,43	*	* * *		runq
