# debian-bo-master/bo/source/net/nis-2.20/debian/postrm
#!/bin/bash
# post removal script for the Debian GNU/Linux nis package

if [ $1 = "purge" ]
then
  update-rc.d nis remove >/dev/null
fi

exit 0
